# march/root/etc/systemd/system/getty@tty1.service.d/autologin.conf
[Unit]
Description=Auto Login
After=systemd-user-sessions.service

[Service]
Type=simple
ExecStart=
ExecStart=-/usr/bin/agetty --autologin march --noclear %I 38400 linux
